# deploy/tunnelflux.service
[Unit]
Description=TunnelFlux Cloudflare edge IP steering for cloudflared
After=network-online.target cloudflared.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tf run --config /etc/tunnelflux/config.json
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
